# macos/launchd/com.marx.agent-inbox.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.marx.agent-inbox</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>__HOME__/.dotfiles/macos/scripts/agent-inbox-daemon.py</string>
  </array>
  <!-- No secrets here: token comes from the Keychain (service
       agent-inbox-token), chat ID from ~/.config/agent-inbox/env -->
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/agent-inbox/daemon.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/agent-inbox/daemon.err</string>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
